# examples/systemd/gosystem-test-mcp.service
[Unit]
Description=Gosystem Test MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/gosystem-test-mcp
Environment=GOSYSTEM_MCP_TRANSPORT=streamable-http
Environment=GOSYSTEM_MCP_HOST=0.0.0.0
Environment=GOSYSTEM_MCP_PORT=8000
Environment=GOSYSTEM_MCP_PATH=/mcp
Environment=GOSYSTEM_MCP_STATELESS_HTTP=true
Environment=GOSYSTEM_MCP_JSON_RESPONSE=true
Environment=GOSYSTEM_MCP_CONFIG_TOML=/opt/gosystem-test-mcp/config.toml
ExecStart=/opt/gosystem-test-mcp/.venv/bin/gosystem-test-mcp
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
